# smf-cacert.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='system/ca-certificates' type='service' version='0'>
    <create_default_instance enabled='true'/>
    <single_instance/>
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/minimal:default'/>
    </dependency>
    <exec_method name='start' type='method' exec='/lib/svc/method/svc-cacert' timeout_seconds='60'/>
    <exec_method name='refresh' type='method' exec='/lib/svc/method/svc-cacert' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec=':true' timeout_seconds='60'/>
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='transient'/>
    </property_group>
    <property_group name='config' type='application'>
      <propval name='capath' type='astring' value='/etc/certs/CA'/>
      <propval name='cabundle' type='astring' value='/etc/certs/ca-certificates.crt'/>
      <propval name='opensslcerts' type='astring' value='/etc/openssl/certs'/>
    </property_group>

    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Root CA Certificates Installer</loctext>
      </common_name>
      <documentation>
        <manpage title='svc-cacert' section='1M' manpath='/usr/man'/>
      </documentation>
    </template>
  </service>
</service_bundle>
